# units/caddyfile-reload.path
[Unit]
Description=Caddyfile change watcher

[Path]
PathChanged=/etc/Caddyfile

[Install]
WantedBy=multi-user.target